# infra/systemd/wecom-room-robot-rag-cache-sync.service
[Unit]
Description=WeCom room robot RAG inventory cache refresh
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/wecom-room-robot
EnvironmentFile=/opt/wecom-room-robot/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/wecom-room-robot/.venv/bin/python /opt/wecom-room-robot/scripts/refresh_rag_inventory_cache.py
User=root
Group=root
Nice=5
TimeoutStartSec=600
