# infra/systemd/wecom-room-robot-callback-watchdog.timer
[Unit]
Description=Run WeCom KF callback watchdog every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true
Unit=wecom-room-robot-callback-watchdog.service

[Install]
WantedBy=timers.target
